# deploy/ansible/roles/mqttrooper/templates/mqttrooper.service.j2
[Unit]
Description=MQTTrooper - MQTT/HTTP command daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ mqttrooper_user }}
Group={{ mqttrooper_group }}
WorkingDirectory={{ mqttrooper_working_dir }}
ExecStart={{ mqttrooper_install_dir }}/mqttrooper -c {{ mqttrooper_config_path }}
Environment=PATH={{ mqttrooper_env_path }}
StandardOutput=append:{{ mqttrooper_config.daemon.log_file_path }}
StandardError=append:{{ mqttrooper_config.daemon.error_file_path }}
Restart=always
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ mqttrooper_log_dir }} {{ mqttrooper_working_dir }}
PrivateTmp=true

[Install]
WantedBy=multi-user.target
